# docker/configs/supervisor/worker-laravel-cron.conf
; supervisor laravel-cron config file

[program:laravel-cron]
process_name=%(program_name)s_%(process_num)02d
command=php -dxdebug.mode=off /var/www/html/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=1000
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/cron.log
stopwaitsecs=1800
priority=3
startretries=30
